# lib/crontab
#
# Crontab addition for Duplicacy Scripts
#

# Backups and Maintenance
#Min Hour DoM Mon DoW  Command
45   0    *   *   *    __BIN__/all-actions
